# crontab
*/5 * * * * /home/jahel/server-bedrock/worlds/Bedrock level/script_git.sh
